# ansible/roles/passbolt-autofill/templates/passbolt-health.service.j2
[Unit]
Description=Passbolt Health Check
After=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/passbolt-health-check
User=root
