# ansible/roles/passbolt-autofill/templates/passbolt-health.timer.j2
[Unit]
Description=Passbolt Health Check Timer

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
